# examples/journalserver/service/journalserver.service
[Unit]
Description=Journal server

[Service]
ExecStart=/usr/bin/journalserver
Restart=always

[Install]
WantedBy=multi-user.target
