# examples/lldp_discoverd/service/lldp_discoverd.service
[Unit]
Description=lldp discoverd service

[Service]
ExecStart=/usr/bin/lldp_discoverd
Restart=always

[Install]
WantedBy=multi-user.target
