# examples/event_broker/service/event_broker.service
[Unit]
Description=event broker

[Service]
ExecStart=/usr/bin/event_broker
Restart=always
Type=dbus

[Install]
WantedBy=multi-user.target
